# ansible/roles/app/templates/piggypulse-stack.service.j2
[Unit]
Description=PiggyPulse production stack
Wants=network-online.target docker.service
After=network-online.target docker.service

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory={{ piggypulse_compose_dir }}
ExecStartPre={{ piggypulse_compose_bin }}{% if piggypulse_compose_subcommand %} {{ piggypulse_compose_subcommand }}{% endif %} pull
ExecStart={{ piggypulse_compose_bin }}{% if piggypulse_compose_subcommand %} {{ piggypulse_compose_subcommand }}{% endif %} up -d --remove-orphans
ExecStop={{ piggypulse_compose_bin }}{% if piggypulse_compose_subcommand %} {{ piggypulse_compose_subcommand }}{% endif %} down
TimeoutStartSec=0

[Install]
WantedBy=multi-user.target
